# gated-downloads/ethii-downloads.service
[Unit]
Description=ETHII gated download service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/ethii-downloads/serve.py
Restart=always
RestartSec=3
User=root
WorkingDirectory=/opt/ethii-downloads

[Install]
WantedBy=multi-user.target
